# ansible/roles/postgresql/templates/postgresql.service.j2
[Unit]
Description=PostgreSQL database server
After=syslog.target
After=network.target

[Service]
Type=forking

User=postgres
Group=postgres

Environment=PGDATA=/var/lib/pgsql/{{ pg_version }}/data/

OOMScoreAdjust=-1000

ExecStartPre=/usr/pgsql-{{ pg_version }}/bin/postgresql95-check-db-dir ${PGDATA}
ExecStart=/usr/pgsql-{{ pg_version }}/bin/pg_ctl start -D ${PGDATA} -s -w -t 120
ExecStop=/usr/pgsql-{{ pg_version }}/bin/pg_ctl stop -D ${PGDATA} -s -m fast
ExecReload=/usr/pgsql-{{ pg_version }}/bin/pg_ctl reload -D ${PGDATA} -s

TimeoutSec=300

[Install]
WantedBy=multi-user.target